# deploy/medical-record-audit.service

---

## ไฟล์ systemd (เวอร์ชันใช้งานจริง)
เลือก **หนึ่ง path**

```ini
# /etc/systemd/system/medical-record-audit.service
[Unit]
Description=Medical Record Audit FE (Next.js)
After=network.target

[Service]
# ปรับผู้ใช้ให้ตรงกับเครื่อง (แนะนำ user ปกติ ไม่ใช่ root)
User=www-data
Group=www-data

# ที่อยู่โปรเจกต์
WorkingDirectory=/opt/medical-record-audit

# โหลดตัวแปรจาก .env (ถ้าไม่มีไฟล์ ให้ข้ามได้ด้วย -)
EnvironmentFile=-/opt/medical-record-audit/.env

# โหมด A: ใช้ Next standalone
ExecStart=/usr/bin/env node .next/standalone/server.js

# โหมด B: ถ้าใช้ custom server.js ให้คอมเมนต์ A แล้วเปิดบรรทัดนี้
# ExecStart=/usr/bin/env node server.js

# ค่าพื้นฐาน
Environment=NODE_ENV=production
Environment=PORT=5108

# restart และทรัพยากร
Restart=on-failure
RestartSec=3
LimitNOFILE=65535

# hardening เบื้องต้น
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
CapabilityBoundingSet=
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true

# logs
SyslogIdentifier=medical-record-audit
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
